# deploy/camdetect.service
[Unit]
Description=camdetect multi-camera detection web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
Environment=CAMDETECT_MODE=__MODE__
Environment=HOST=0.0.0.0
Environment=PORT=__PORT__
ExecStart=__WORKDIR__/run.sh __MODE__
# Auto-recovery: restart on any crash, keep retrying.
Restart=always
RestartSec=5
StartLimitIntervalSec=0
KillMode=mixed
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
